# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
loglevel=info
user=root

[program:postgresql]
command=/bin/bash -c "su postgres -c 'postgres -D /var/lib/postgresql/data -p 5432'"
stdout_logfile=/var/log/postgresql.log
stderr_logfile=/var/log/postgresql.log
autorestart=true
startretries=5
priority=10

[program:backend]
command=node /app/server.cjs
directory=/app
stdout_logfile=/var/log/backend.log
stderr_logfile=/var/log/backend.log
autorestart=true
startretries=10
startsecs=3
priority=20
environment=PORT="3001",DATABASE_URL="%(ENV_DATABASE_URL)s",JWT_SECRET="%(ENV_JWT_SECRET)s"

[program:frontend]
command=node /frontend/server.js
directory=/frontend
stdout_logfile=/var/log/frontend.log
stderr_logfile=/var/log/frontend.log
autorestart=true
startretries=5
startsecs=3
priority=20
environment=PORT="3000",HOSTNAME="0.0.0.0",NODE_ENV="production"

[program:nginx]
command=nginx -g 'daemon off;'
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx.log
autorestart=true
startretries=3
priority=30
